# templates/tmpws.service.j2
[Unit]
After=network.target
[Service]
Type=simple
WorkingDirectory=/opt
ExecStart=/usr/bin/python -m SimpleHTTPServer 1234
[Install]
WantedBy=default.target